# services/mud/run
#!/bin/bash
cd /app/mud
node /app/mud/app.js